# build/crontab
# LOG HONEYWELL THERMOSTAT DATA EVERY 5 MIN
*/5 * * * * /scripts/honeywell-logger.sh | sed '/^\s*$/d' | awk '{ print strftime("[%Y-%m-%d %H:%M:%S]"), $0 }' >>/var/log/honeywell-logger/honeywell-logger.log 2>&1
*/5 * * * * /scripts/weather-logger.sh | sed '/^\s*$/d' | awk '{ print strftime("[%Y-%m-%d %H:%M:%S]"), $0 }' >>/var/log/honeywell-logger/honeywell-logger.log 2>&1
