# build/crontab
# LOG HONEYWELL THERMOSTAT DATA EVERY 5 MIN
*/5 * * * * /scripts/honeywell-logger.sh | ts '%Y-%m-%d %H:%M:%S %Z -' >>/var/log/honeywell-logger/honeywell-logger.log 2>&1
*/5 * * * * /scripts/weather-logger.sh | ts '%Y-%m-%d %H:%M:%S %Z -' >>/var/log/honeywell-logger/honeywell-logger.log 2>&1
# BLANK LINE AT END OF CRONTAB
